# .devcontainer/s6-rc.d/docker-entrypoint/up
#!/command/execlineb -P
/usr/local/bin/docker-entrypoint.sh
