# .devcontainer/s6-rc.d/supervisord/run
#!/command/with-contenv bash
# s6-overlay の with-contenv で環境変数を読み込む
set -euo pipefail
exec /usr/bin/supervisord -c /etc/supervisor/supervisord.conf