# .devcontainer/s6-rc.d/process-compose/run
#!/usr/bin/env bash
set -euo pipefail
exec /usr/local/bin/process-compose -f /etc/process-compose/process-compose.yaml